# root/etc/services.d/lazylibrarian/run
#!/usr/bin/with-contenv bash

exec \
	s6-setuidgid abc python3 /app/lazylibrarian/LazyLibrarian.py \
	--datadir /config --nolaunch
